# backend/crontab
0 2 * * * cd /app && python manage.py backup_to_telegram >> /proc/1/fd/1 2>/proc/1/fd/2
